# resources/config/supervisord.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/data1/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:es]
command=/usr/local/elasticsearch/bin/elasticsearch
directory=/usr/local/elasticsearch
autostart=true
autorestart=true
startretries=10
stopsignal=INT
stopwaitsecs=300
stderr_logfile=/data1/logs/supervisor_es_error.log
stdout_logfile=/data1/logs/supervisor_es.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
user=ec2-user
environment=JAVA_HOME=/usr/java/default,HOME='/home/ec2-user'


[program:gunicorn]
command=/usr/local/bin/gunicorn --config=resources/config/gunicorn.py active_data.app:flask_app
directory=/home/ec2-user/ActiveData
autostart=true
autorestart=true
startretries=10
stderr_logfile=/data1/logs/supervisor_gunicorn_error.log
stdout_logfile=/data1/logs/supervisor_gunicorn.log
user=ec2-user
environment=PYTHONPATH='.:vendor',PYPY_GC_MAX='6GB',HOME='/home/ec2-user',ACTIVEDATA_CONFIG=resources/config/staging_settings.json

[program:tuid]
command=/user/bin/python tuid/app.py --config=TUID/resources/config/prod.json
directory=/home/ec2-user/TUID
autostart=true
autorestart=true
startretries=10
stderr_logfile=/data1/logs/supervisor_tuid_error.log
stdout_logfile=/data1/logs/supervisor_tuid.log
user=ec2-user
environment=PYTHONPATH='.:vendor',PYPY_GC_MAX='6GB',HOME='/home/ec2-user'
