# resources/config/supervisord.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:es]
command=/usr/local/elasticsearch/bin/elasticsearch
directory=/usr/local/elasticsearch
autostart=true
autorestart=true
startretries=10
stopsignal=INT
stopwaitsecs=60
stderr_logfile=/logs/es.error.log
stdout_logfile=/logs/es.log
user=root
environment=ES_MIN_MEM=5g;ES_MAX_MEM=5g

[program:flask]
command=python27 active_data/app.py --settings=resources/config/staging_settings.json --process=%(process_num)d
directory=/home/ec2-user/ActiveData
numprocs=2
process_name=%(process_num)s
autostart=true
autorestart=true
startretries=10
stderr_logfile=/logs/flask%(process_num)s.error.log
stdout_logfile=/logs/flask%(process_num)s.log
user=ec2-user
environment=PYTHONPATH='.',HOME='/home/ec2-user'

# [program:gunicorn]
# command=/usr/local/bin/gunicorn --pythonpath . --config resources/config/gunicorn.py 'active_data.app:setup(settings="resources/config/staging_settings.json")'
# directory=/home/ec2-user/ActiveData
# autostart=true
# autorestart=true
# startretries=10
# stderr_logfile=/logs/gunicorn.error.log
# stdout_logfile=/logs/gunicorn.log
# user=ec2-user
# environment=PYTHONPATH='.',PYPY_GC_MAX='6GB',HOME='/home/ec2-user'
