# resources/docker/supervisord.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:gunicorn]
command=/usr/local/bin/gunicorn --config=resources/config/gunicorn.py active_data.app:flask_app
directory=/app
autostart=true
autorestart=true
startretries=10
stderr_logfile=/app/logs/supervisor_gunicorn_error.log
stdout_logfile=/app/logs/supervisor_gunicorn.log
user=root
environment=PYTHONPATH='.:vendor',HOME='/app',ACTIVEDATA_CONFIG=resources/config/staging_settings.json

